# scripts/dispatcher-worker.service
[Unit]
Description=dispatcher-worker.service

[Service]
Environment=PYTHONPATH=/root/scanning
ExecStart=/root/scanning/workers/dispatcher.py
Restart=on-failure
