# scripts/notify-worker.service
[Unit]
Description=notify-worker.service

[Service]
Environment=PYTHONPATH=/root/analytics-integration
ExecStart=/root/analytics-integration/mail_service/notify_worker.py
Restart=on-failure
