# scripts/scan-worker.service
[Unit]
Description=scan-worker.service

[Service]
Environment=PYTHONPATH=/root/scanning
ExecStart=/root/scanning/workers/scan.py
Restart=on-failure
