# scripts/scan-worker.service
[Unit]
Description=scan-worker.service

[Service]
Environment=PYTHONPATH=/opt/scanning
ExecStart=/opt/scanning/workers/scan.py
Restart=on-failure
